# deployment/files/init/mozio-web.conf
start on starting mozio
stop on stopping mozio